# .config/systemd/user/dst-discord-bot.service
[Unit]
Description=DST Discord Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/dst-fish-manager
ExecStart=/usr/bin/python3 %h/dst-fish-manager/run_discord_bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=default.target
