# config/systemd/user/dontstarve.target
[Unit]
Description=Don't Starve Together dedicated server target (group target)

[Install]
WantedBy=default.target
